# .deploy/supervisor/conf.d/lawyerd.conf
#[unix_http_server]
#file=/var/run/supervisor.sock   ; (the path to the socket file)
#chmod=0700 ; sockef file mode (default 0700)
#username = gohrmmhdfgggrtg45g
#password = g34g42fg34g45g45grht

[supervisord]
logfile = /var/log/supervisord_lawyerd.log
logfile_maxbytes = 50MB
logfile_backups=10
; loglevel = warn
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
;directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false
;user=www-data
user=root
environment=DJANGO_READ_DOT_ENV_FILE="True"


#[supervisorctl]
#serverurl=unix:///var/run//supervisor.sock ; use a unix:// URL for a unix socket

[program:gunicorn]
command=/home/lawyerd/venv/bin/gunicorn config.wsgi -b 127.0.0.1:8000 --chdir=/home/lawyerd/ -e DJANGO_SETTINGS_MODULE="config.settings.production" -e DJANGO_READ_DOT_ENV_FILE="True" --error-logfile /var/log/error_file_g.log --log-level=debug
;command=/home/lawyerd/venv/bin/python3.7 /home/lawyerd/manage.py collectstatic

;command=/home/lawyerd/venv/bin/gunicorn config.wsgi:application -b unix:/tmp/gunicorn.sock -e DJANGO_READ_DOT_ENV_FILE="True"
;command=/home/lawyerd/venv/bin/gunicorn config.wsgi:application -b 0.0.0.0:8081 -e DJANGO_READ_DOT_ENV_FILE="True"
;,DJANGO_DEBUG=False,DJANGO_SETTINGS_MODULE=config.settings.production
;directory=/home/lawyerd/
stdout_logfile=/var/log/lawyerd_out.log
;stderr_logfile=/var/log/lawyerd_error.log

;user=nobody
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
;environment=DJANGO_READ_DOT_ENV_FILE=True,DJANGO_DEBUG=False
startsecs=10
stopasgroup=true



[program:celery_worker]
#command=export CELERY_BROKER_URL=redis://localhost:6379/1 && export DJANGO_READ_DOT_ENV_FILE=True && /home/lawyerd/venv/bin/python -m celery -A config.celery_app worker -P solo --loglevel=INFO --workdir /home/lawyerd/lawyerd/
command=/home/lawyerd/venv/bin/python3.7 -m celery -A config.celery_app worker -P solo --loglevel=INFO --without-gossip --without-mingle --without-heartbeat -Ofair --workdir /home/lawyerd/lawyerd
directory=/home/lawyerd/
#command=/home/venv/bin/python3.7 /home/autosubscriber/manage.py domain_subscribe
autostart=true
autorestart=true
stdout_logfile=/var/log/celery_worker_out.log
stderr_logfile=/var/log/celery_worker_error.log
numprocs=1
#process_name = %(program_name)s_%(process_num)02d
stopsignal=QUIT
#environment=PYTHONUNBUFFERED=1,DJANGO_SETTINGS_MODULE="config.settings.local", DATABASE_URL="postgres://postgres:8hg584g8b3b3bf49H9bvbGUIGUIkbev33@localhost:5432/autosubscriber", CELERY_BROKER_URL="redis://localhost:6379/0", CELERY_RESULT_BACKEND="redis://localhost:6379/0", USE_DOCKER="False", SENTRY_DSN="https://b566074a8281474abcaeb3bb7fed3143@sentry.io/1519321"


startsecs=10
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600
; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true
; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000


[program:celery_beat]
command=/home/lawyerd/venv/bin/python3.7 -m celery -A config.celery_app beat -l info --pidfile= --workdir /home/lawyerd/lawyerd
directory=/home/lawyerd/
#command=/home/venv/bin/python3.7 /home/autosubscriber/manage.py domain_subscribe
autostart=true
autorestart=true
stdout_logfile=/var/log/celery_beat_out_lawyerd.log
stderr_logfile=/var/log/celery_beat_error_lawyerd.log
numprocs=1
#process_name = %(program_name)s_%(process_num)02d
stopsignal=QUIT


startsecs=10
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600
; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true
; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000




#[Tprogram:local_site]
#command=/home/lawyerd/venv/bin/python3.7 /home/lawyerd/manage.py runserver lawyerd.net:8000 --settings=config.settings.local
#command=/home/lawyerd/venv/bin/python3.7 /home/lawyerd/manage.py runserver 167.99.243.12:8000 --settings=config.settings.local
#command=/home/lawyerd/venv/bin/python3.7 /home/lawyerd/manage.py runserver 167.99.243.12:8000 --settings=config.settings.local
#autostart=true
#autorestart=true
#stdout_logfile=/var/log/site_out_lawyerd.log
#stderr_logfile=/var/log/site_error_lawyerd.log

#numprocs=1
#process_name = %(program_name)s_%(process_num)02d


#priority=996
stopsignal=QUIT
